# qrtr-ns.service
[Unit]
Description=QIPCRTR Name Service

[Service]
ExecStart=/usr/local/bin/qrtr-ns -f 1
Restart=always

[Install]
WantedBy=multi-user.target
